# provision/bigrock-update.service
[Unit]
Description=BigRock Updater

[Service]
Type=oneshot
User=pi
ExecStart=/opt/bigrock-app/scripts/update.sh
